# deploy/logrotate/gmail-search-logrotate.timer
[Unit]
Description=Daily gmail-search log rotation

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
